# systemd/contabo-storage-python.service
[Unit]
Description=Contabo Storage Manager - Python (FastAPI)
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/contabo_storage_manager/packages/python-bridge
EnvironmentFile=/root/contabo_storage_manager/.env
ExecStart=/usr/bin/python3 -m uvicorn app.main:app --host 0.0.0.0 --port 8000 --workers 2
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=contabo-storage-python

[Install]
WantedBy=multi-user.target
